# deploy/macos/com.secretar.compose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd — автозапуск Secretar (macOS).
     Установка:
       sed -i '' 's|__REPO__|/полный/путь/к/secretar-repo|' com.secretar.compose.plist
       cp com.secretar.compose.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.secretar.compose.plist
     Docker Desktop: Settings -> Start when you log in. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.secretar.compose</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>sleep 120 &amp;&amp; cd __REPO__/infra &amp;&amp; /usr/local/bin/docker compose up -d</string>
  </array>
  <key>RunAtLoad</key><true/>
</dict>
</plist>
